# .config/systemd/user/pacaur-sync.service
[Unit]
Description=Update local database of package manager.

[Service]
Type=oneshot
ExecStart=pacaur -Sy
RemainAfterExit=true

[Install]
WantedBy=default.target
